# crontab
* * * * * /Users/timborn/pinggraph/bin/zping.sh > /dev/null 2>&1
# * * * * * /Users/timborn/pinggraph/bin/zping.sh 
